# deploy/alunotes-web.service
[Unit]
Description=AluNotes Bridge Web Control Plane
After=network.target alunotes-bridge.service alunotes-ai.service
Wants=alunotes-bridge.service alunotes-ai.service

[Service]
Type=simple
User=@SERVICE_USER@
Environment=NODE_ENV=production
Environment=XDG_RUNTIME_DIR=/run/user/@SERVICE_UID@
Environment=HOME=/home/@SERVICE_USER@
WorkingDirectory=@PROJECT_DIR@/alunotes-bt-web
# Runtime lives in the distrobox because Node/pnpm are not installed on the host.
ExecStart=/usr/bin/distrobox enter @DISTROBOX_NAME@ -- sh -lc 'cd @PROJECT_DIR@/alunotes-bt-web && exec pnpm start'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
